# src/bk/data/systemd/bk-worker.service
[Unit]
Description=GPUbk per-user scheduled job worker
After=default.target

[Service]
Type=simple
EnvironmentFile=-%h/.config/bk/bk.env
ExecStart=@PYTHON_EXECUTABLE@ -m bk worker
Restart=on-failure
RestartSec=5
TimeoutStopSec=15
KillMode=control-group
NoNewPrivileges=true

[Install]
WantedBy=default.target
