# src/bk/data/systemd/bk-worker.service
[Unit]
Description=GPUbk per-user scheduled job worker
After=default.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
@SERVICE_ENVIRONMENT@
ExecStart=@PYTHON_EXECUTABLE@ -m bk worker
Restart=on-failure
RestartPreventExitStatus=75
RestartSec=5
TimeoutStopSec=15
KillMode=control-group
NoNewPrivileges=true

[Install]
WantedBy=default.target
