# deploy/bk-monitor.service
[Unit]
Description=BK GPU usage monitor
After=default.target

[Service]
Type=simple
EnvironmentFile=-%h/.config/bk/bk.env
ExecStart=%h/.local/bin/bk monitor --interval 2 --rollup 60
Restart=on-failure
RestartSec=5
Nice=10
IOSchedulingClass=idle
NoNewPrivileges=true

[Install]
WantedBy=default.target
